# udev/99-network-cleanup.rules
SUBSYSTEM=="net", ACTION=="remove", RUN+="/usr/sbin/network-cleanup.sh $env{INTERFACE}"
